# evil_nmap.service
[Unit]
Description=EvilCrow RF Scanner Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/evilcrow
ExecStart=/opt/evilcrow/venv/bin/python /opt/evilcrow/evilcrow_scanner.py --interval 900 --ip 192.168.4.1 --freqs 433,915,315,868,405
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=evilcrow

[Install]
WantedBy=multi-user.target
